# init/kronosnetd.service.in
[Unit]
Description=kronosnetd
Requires=network.target
After=network.target syslog.target

[Service]
Type=forking
EnvironmentFile=/etc/sysconfig/kronosnetd
ExecStart=@SBINDIR@/kronosnetd

[Install]
WantedBy=multi-user.target
