# knetd/systemd/knetd@.service
[Unit]
Description=Kronosnet VPN Daemon (Instance: %i)
Documentation=man:knetd(8) man:knetd.conf(5)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/sbin/knetd --config /etc/knetd/%i.conf --node-id ${NODE_ID}
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s

# Security hardening
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/knetd /run/knetd /var/log/knetd
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Resource limits
LimitNOFILE=65536
LimitMEMLOCK=infinity

# State directory
StateDirectory=knetd/%i
RuntimeDirectory=knetd/%i
LogsDirectory=knetd/%i
ConfigurationDirectory=knetd

# Environment - NODE_ID must be set via systemctl edit or drop-in
EnvironmentFile=-/etc/knetd/%i.env

[Install]
WantedBy=multi-user.target
